# weston.service
[Unit]
Description=Weston Wayland Compositor (Headless)
Documentation=man:weston(1)

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/workspace
ExecStart=/usr/bin/weston \
  --backend=headless-backend.so \
  --xwayland \
  --width=1920 \
  --height=1080 \
  --socket=wayland-0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target